# test-apm.service
# pod-31e1b6136b67967c59bb1e3096b5c8ce1880d0b2a7ca62f9467123d22221bbf4.service
# autogenerated by Podman 1.6.4
# Tue Jun 28 16:12:47 KST 2022

[Unit]
Description=Podman pod-31e1b6136b67967c59bb1e3096b5c8ce1880d0b2a7ca62f9467123d22221bbf4.service
Documentation=man:podman-generate-systemd(1)
Requires=container-39f98a2aa1f0ea3819665523c7bdfd07fe04bc5d44178f3ec18711dbdd7fcb34.service container-bdf03e02eb8d5cbb0c7c20d7913abcc2934bd95547c9a5bb02681fe5c826ba4a.service container-e2114796cda330411818259b56a95d192bdab9c418bd6442c10d23477d278cef.service
Before=container-39f98a2aa1f0ea3819665523c7bdfd07fe04bc5d44178f3ec18711dbdd7fcb34.service container-bdf03e02eb8d5cbb0c7c20d7913abcc2934bd95547c9a5bb02681fe5c826ba4a.service container-e2114796cda330411818259b56a95d192bdab9c418bd6442c10d23477d278cef.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start 012c9a735ecd3221064b8ad011018f20de956494804b523501dbc6f9e5f5b4a7
ExecStop=/usr/bin/podman stop -t 10 012c9a735ecd3221064b8ad011018f20de956494804b523501dbc6f9e5f5b4a7
KillMode=none
Type=forking
PIDFile=/var/run/containers/storage/overlay-containers/012c9a735ecd3221064b8ad011018f20de956494804b523501dbc6f9e5f5b4a7/userdata/conmon.pid

[Install]
WantedBy=multi-user.target
# container-39f98a2aa1f0ea3819665523c7bdfd07fe04bc5d44178f3ec18711dbdd7fcb34.service
# autogenerated by Podman 1.6.4
# Tue Jun 28 16:12:47 KST 2022

[Unit]
Description=Podman container-39f98a2aa1f0ea3819665523c7bdfd07fe04bc5d44178f3ec18711dbdd7fcb34.service
Documentation=man:podman-generate-systemd(1)
RefuseManualStart=yes
RefuseManualStop=yes
BindsTo=pod-31e1b6136b67967c59bb1e3096b5c8ce1880d0b2a7ca62f9467123d22221bbf4.service
After=pod-31e1b6136b67967c59bb1e3096b5c8ce1880d0b2a7ca62f9467123d22221bbf4.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start 39f98a2aa1f0ea3819665523c7bdfd07fe04bc5d44178f3ec18711dbdd7fcb34
ExecStop=/usr/bin/podman stop -t 10 39f98a2aa1f0ea3819665523c7bdfd07fe04bc5d44178f3ec18711dbdd7fcb34
KillMode=none
Type=forking
PIDFile=/var/run/containers/storage/overlay-containers/39f98a2aa1f0ea3819665523c7bdfd07fe04bc5d44178f3ec18711dbdd7fcb34/userdata/conmon.pid

[Install]
WantedBy=multi-user.target
# container-bdf03e02eb8d5cbb0c7c20d7913abcc2934bd95547c9a5bb02681fe5c826ba4a.service
# autogenerated by Podman 1.6.4
# Tue Jun 28 16:12:47 KST 2022

[Unit]
Description=Podman container-bdf03e02eb8d5cbb0c7c20d7913abcc2934bd95547c9a5bb02681fe5c826ba4a.service
Documentation=man:podman-generate-systemd(1)
RefuseManualStart=yes
RefuseManualStop=yes
BindsTo=pod-31e1b6136b67967c59bb1e3096b5c8ce1880d0b2a7ca62f9467123d22221bbf4.service
After=pod-31e1b6136b67967c59bb1e3096b5c8ce1880d0b2a7ca62f9467123d22221bbf4.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start bdf03e02eb8d5cbb0c7c20d7913abcc2934bd95547c9a5bb02681fe5c826ba4a
ExecStop=/usr/bin/podman stop -t 10 bdf03e02eb8d5cbb0c7c20d7913abcc2934bd95547c9a5bb02681fe5c826ba4a
KillMode=none
Type=forking
PIDFile=/var/run/containers/storage/overlay-containers/bdf03e02eb8d5cbb0c7c20d7913abcc2934bd95547c9a5bb02681fe5c826ba4a/userdata/conmon.pid

[Install]
WantedBy=multi-user.target
# container-e2114796cda330411818259b56a95d192bdab9c418bd6442c10d23477d278cef.service
# autogenerated by Podman 1.6.4
# Tue Jun 28 16:12:47 KST 2022

[Unit]
Description=Podman container-e2114796cda330411818259b56a95d192bdab9c418bd6442c10d23477d278cef.service
Documentation=man:podman-generate-systemd(1)
RefuseManualStart=yes
RefuseManualStop=yes
BindsTo=pod-31e1b6136b67967c59bb1e3096b5c8ce1880d0b2a7ca62f9467123d22221bbf4.service
After=pod-31e1b6136b67967c59bb1e3096b5c8ce1880d0b2a7ca62f9467123d22221bbf4.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/podman start e2114796cda330411818259b56a95d192bdab9c418bd6442c10d23477d278cef
ExecStop=/usr/bin/podman stop -t 10 e2114796cda330411818259b56a95d192bdab9c418bd6442c10d23477d278cef
KillMode=none
Type=forking
PIDFile=/var/run/containers/storage/overlay-containers/e2114796cda330411818259b56a95d192bdab9c418bd6442c10d23477d278cef/userdata/conmon.pid

[Install]
WantedBy=multi-user.target
